# devops/ansible/exemples/roles/mysql-exporter/templates/mysqld_exporter.service.j2
[Unit]
Description=MySQL Node Exporter
After=network.target
 
[Service]
Type=simple
ExecStart=/usr/local/bin/mysqld_exporter
 
[Install]
WantedBy=multi-user.target
